# systemd/acpi-volumed.service
[Unit]
Description=ACPI Volume Daemon
#Requires=acpid.service

[Service]
ExecStart=/usr/bin/acpi-volumed

[Install]
WantedBy=default.target
